# network/if-up.d/add-my-static-arp
arp -i eth0 -s 191.193.177.254 ca:fe:c0:ff:ee:00
